# debian/libra_statsd.upstart
description "statsd"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

script
   libra_statsd -c /etc/libra.cfg 

end script

