# debian/libra_admin_api.upstart
description "libra_admin_api"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

script
   libra_admin_api -c /etc/libra.cfg

end script

